# CourseWorks/cw3/tvs-2425-1_cw3/tvsctl-srv/etc/service/tvsctld.service
[Unit]
Description=TVS Control Application

[Service]
ExecStart=/opt/isel/tvs/tvsctld/bin/tvsctld
Restart=no

[Install]
WantedBy=multi-user.target

